# sniproxy.service

[Unit]
Description= sni proxy
After=network.target

[Service]
ExecStart=/usr/sbin/sniproxy -f
Restart=always

[Install]
WantedBy=multi-user.target
